# notifier.service
[Unit]
Description=Periodic desktop notification daemon
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/notifier __daemon__
Restart=on-failure
RestartSec=2
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/%U/bus

[Install]
WantedBy=default.target
